# dot_config/systemd/user/push-etckeeper.timer
[Unit]
Description=Daily git push etckeeper repo

[Timer]
OnStartupSec=1m
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target 
