# systemd/user/offlineimap-oneshot.timer.d/override.conf
[Timer]
OnBootSec=
OnStartupSec=5m
OnUnitInactiveSec=
OnUnitInactiveSec=2h
