# systemd/user/vdirsyncer.timer.d/override.conf
[Timer]
OnBootSec=
OnUnitActiveSec=
AccuracySec=
OnStartupSec=5m
OnUnitActiveSec=2h
