# systemd/user/vdirsyncer.timer.d/override.conf
[Timer]
OnBootSec=
OnUnitActiveSec=
AccuracySec=
OnBootSec=30s
OnUnitActiveSec=1h
