# systemd/user/backup-bootloader-settings.service
[Unit]
Description=Copy bootloader configs to ~/.config

[Service]
Type=oneshot
ExecStart=/bin/sh -c '/usr/bin/rsync -a /boot/loader ~/.config/boot'
