# systemd/user/backup-bootloader-settings.timer
[Unit]
Description=Run bootloader backup after boot and once a week

[Timer]
OnBootSec=30s
OnUnitActiveSec=1w

[Install]
WantedBy=timers.target 
