# systemd/user/backup-etc-files.timer
[Unit]
Description=Run etc-folder backup after boot and once a day

[Timer]
OnStartupSec=5m
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target 
