# systemd/user/backup-bootloader-settings.timer
[Unit]
Description=Run bootloader backup after boot and every week

[Timer]
OnStartupSec=1m
OnUnitActiveSec=1w

[Install]
WantedBy=timers.target 
