# systemd/user/push-etckeeper.service
[Unit]
Description=Push etckeeper repo to remote

[Service]
Type=oneshot
ExecStart="/home/otahontas/.config/systemd/user/scripts/push-etckeeper.sh"
